# cpufreq.service
[Unit]
Description=CPU powersave
[Service]
Type=oneshot
ExecStart=/usr/bin/cpufreq-setup-PE 400 4200
[Install]
WantedBy=multi-user.target
